# wh2900-monitor.timer
[Unit]
Description=WH2900 Station Online Monitor Timer

[Timer]
OnBootSec=300
OnUnitActiveSec=900
AccuracySec=60

[Install]
WantedBy=timers.target
